# systemd-user/archive-notes.service
[Unit]
Description=Archive past-month daily notes and logs

[Service]
Type=oneshot
ExecStart=/home/pratik/scripts/archive-notes.sh
