# systemd-user/archive-notes.timer
[Unit]
Description=Monthly archive of daily notes and logs

[Timer]
OnCalendar=*-*-01 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
