# systemd-user/cram.service
[Unit]
Description=OS cram — tick notifications

[Service]
Type=oneshot
ExecStart=%h/scripts/cram tick
