# systemd-user/disk-cleanup.service
[Unit]
Description=Periodic disk cleanup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/disk-cleanup.sh
