# systemd-user/dotfiles-backup.service
[Unit]
Description=Backup dotfiles to GitHub

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-backup.sh --auto
