# systemd-user/habits-prompt-meditate-0900.timer
[Unit]
Description=Habit prompt: meditate @ 09:00

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=false
AccuracySec=30s
Unit=habits-prompt@meditate_0900.service

[Install]
WantedBy=timers.target
