# systemd-user/habits-catchup.service
[Unit]
Description=Habits catchup (mark missed slots)

[Service]
Type=oneshot
ExecStart=/home/pratik/.local/bin/habits catchup
