# systemd-user/hypr-backup.service
[Unit]
Description=Backup Hyprland config to GitHub

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hypr-backup.sh
